# files/etc/cron.daily/ntpd
#!/bin/sh
/usr/bin/ntpdate pool.ntp.org


